# raspberry-pos-os/stage-pos/04-pos-system/files/pos-chain-verify.timer
[Unit]
Description=Dagelijkse SBA hash-chain verificatie om 04:00

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
